# supervisor/prod.conf
[program:webserver]
directory = /opt/ca
command = gunicorn project.wsgi --bind=0.0.0.0:80 --workers=5
stdout_logfile = /var/log/ca/webserver.log
stderr_logfile = /var/log/ca/webserver.err